# yyy.conf
[program:yyy]
command=python3 /home/curator2/yyy/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/yyy/test.err.log
stdout_logfile=/home/curator2/yyy/test.out.log